# ch_collections/security_banner/roles/banner/files/etc/profile.d/banner.sh
#! /usr/bin/bash

if [ -r /etc/motd ]
then
	cat /etc/motd
	printf '%s\n\n' "Connected to $(hostname)"
fi